# data/switcheroo-control.service.in
[Unit]
Description=Switcheroo Control Proxy service

[Service]
Type=dbus
BusName=net.hadess.SwitcherooControl
ExecStart=@libexecdir@/switcheroo-control

[Install]
WantedBy=graphical.target
