# data/switcheroo-control.service.in
[Unit]
Description=Switcheroo Control Proxy service
Before=multi-user.target display-manager.target alsa-restore.service alsa-state.service

[Service]
Type=dbus
BusName=net.hadess.SwitcherooControl
ExecStart=@libexecdir@/switcheroo-control

[Install]
WantedBy=graphical.target
